# startup.service
[Unit]
Description=Startup

[Service]
ExecStart=/usr/bin/startup

[Install]
WantedBy=multi-user.target
